# dist/com.rtw88.client.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.rtw88.client</string>

	<!-- Bring up the RTL8821CE and hand the data path to the kext (enX), reading
	     settings (networks, band, mode) from the config file. Runs as root: it needs
	     the RTW88Server user client plus ifconfig/route/scutil. -->
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/RTW88Client</string>
		<string>--config</string>
		<string>/usr/local/etc/rtw88.conf</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<!-- Keep the link up: if RTW88Client exits (link lost, error, or the kext
	     wasn't ready yet at boot) launchd restarts it -> auto-reconnect. -->
	<key>KeepAlive</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>10</integer>

	<key>StandardOutPath</key>
	<string>/var/log/rtw88.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/rtw88.log</string>
</dict>
</plist>
