# dist/com.rtw88.menu.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.rtw88.menu</string>

	<!-- LaunchAgent for the RTW88 menu-bar app. Runs in the console user's GUI
	     (Aqua) session and starts at login. NOT KeepAlive: the menu has a Quit
	     item, and force-relaunching a UI app the user explicitly quit is hostile —
	     it comes back at the next login via RunAtLoad. The app is pure UI; it talks
	     to the root rtwd daemon over the Unix socket. -->
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/RTW88Menu.app/Contents/MacOS/RTW88Menu</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<!-- only load in a GUI login session (a menu-bar app needs the Aqua session) -->
	<key>LimitLoadToSessionType</key>
	<string>Aqua</string>
</dict>
</plist>
